# deploy/com.chuzom.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- REFERENCE TEMPLATE — do not use verbatim. The __PYTHON__ / __HOME__ tokens are
     placeholders. Generate a copy with your own interpreter/home paths via:
         python -m chuzom.gateway_service
     which writes ~/Library/LaunchAgents/com.chuzom.gateway.plist for this user.
     Then: launchctl load ~/Library/LaunchAgents/com.chuzom.gateway.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chuzom.gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>chuzom.gateway</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CHUZOM_GATEWAY_HOST</key><string>127.0.0.1</string>
        <key>CHUZOM_GATEWAY_PORT</key><string>17900</string>
    </dict>
    <key>KeepAlive</key><true/>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>__HOME__/.chuzom/gateway.out.log</string>
    <key>StandardErrorPath</key><string>__HOME__/.chuzom/gateway.err.log</string>
</dict>
</plist>
