# deploy/com.chuzom.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd agent: keep the Chuzom OpenAI-compatible gateway always running, so any
     external agent (Stockagent, cron jobs) that points OPENAI_BASE_URL at it routes
     through Chuzom. Install: cp to ~/Library/LaunchAgents/ && launchctl load it. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chuzom.gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yaliandrona/Projects/Chuzom/.venv/bin/python</string>
        <string>-m</string>
        <string>chuzom.gateway</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CHUZOM_GATEWAY_HOST</key><string>127.0.0.1</string>
        <key>CHUZOM_GATEWAY_PORT</key><string>17900</string>
    </dict>
    <key>KeepAlive</key><true/>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/Users/yaliandrona/.chuzom/gateway.out.log</string>
    <key>StandardErrorPath</key><string>/Users/yaliandrona/.chuzom/gateway.err.log</string>
</dict>
</plist>
